# ansible/config_files/crontab
SHELL=/bin/bash

* */1 * * * backup sh /home/backup/script_backup.sh